# testsuite/usr/lib/system/libxpc.dylib/resources/launchd_basic_mach_service_template.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Program</key>
	<string>@PROGRAM@</string>
	<key>StandardOutPath</key>
	<string>@STANDARDOUTPATH@</string>
	<key>StandardErrorPath</key>
	<string>@STANDARDERRORPATH@</string>
	<key>Label</key>
	<string>@SERVICE_NAME@</string>
	<key>MachServices</key>
	<dict>
		<key>@SERVICE_NAME@</key>
		<true/>
	</dict>
</dict>
</plist>
